# ops/systemd/watchdirs-vacuum.timer
[Unit]
Description=Run watchdirs weekly SQLite maintenance

[Timer]
OnCalendar=Sun *-*-* 03:17:00
Persistent=true
Unit=watchdirs-vacuum.service

[Install]
WantedBy=timers.target
