# ops/systemd/watchdirs-collect.service
[Unit]
Description=watchdirs hourly evidence collection

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/watchdirs collect --config /etc/watchdirs/watchdirs.toml --db /var/lib/watchdirs/watchdirs.sqlite3 --json --verbose
StateDirectory=watchdirs
ConfigurationDirectory=watchdirs
UMask=0077
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
