# ops/systemd/watchdirs-prune.service
[Unit]
Description=watchdirs retention pruning

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/watchdirs prune --db /var/lib/watchdirs/watchdirs.sqlite3 --hourly-days 3 --daily-days 90 --json
StateDirectory=watchdirs
ConfigurationDirectory=watchdirs
UMask=0077
Nice=19
CPUSchedulingPolicy=idle
CPUWeight=idle
IOSchedulingClass=idle
IOWeight=1
