# ops/systemd/watchdirs-collect.timer
[Unit]
Description=Run watchdirs hourly evidence collection

[Timer]
OnCalendar=hourly
Persistent=true
Unit=watchdirs-collect.service

[Install]
WantedBy=timers.target
